# contour-update.service
[Unit]
Description=Daily auto-update for glass (Docker)

[Service]
Type=oneshot
WorkingDirectory=/opt/glass
ExecStart=/opt/glass/update.sh

